# root/etc/services.d/webdav/run
#!/usr/bin/with-contenv bash

exec \
	lighttpd -f /config/lighttpd.conf